# astra-vein-agent.service
[Unit]
Description=Astra Vein Receptionist Agent
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/roboai/roboai-espeak
Environment="PATH=/home/unitree/.local/bin:/home/unitree/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/unitree/.cargo/bin/uv run src/run.py astra_vein_receptionist

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/unitree/roboai/roboai-espeak/logs/agent.log
StandardError=append:/home/unitree/roboai/roboai-espeak/logs/agent-error.log

# Resource limits (optional - adjust as needed)
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
