# lex_agent_robust.service
[Unit]
Description=Lex Channel Chief Agent - AI Sales Agent for Lexful (Robust)
Documentation=https://github.com/feraco/roboai-espeak
After=network-online.target ollama.service sound.target
Wants=network-online.target ollama.service

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/roboai-espeak
Environment="PATH=/home/unitree/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/unitree/.Xauthority"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="SKIP_AUDIO_VALIDATION=true"

# Wait for system to be fully ready
ExecStartPre=/bin/sleep 10

# Run comprehensive pre-start checks (waits for hardware, validates Ollama)
ExecStartPre=/bin/bash /home/unitree/roboai-espeak/lex_pre_start_checks_robust.sh

# Clear stale audio configs (force re-detection)
ExecStartPre=/bin/bash -c "cd /home/unitree/roboai-espeak && rm -f device_config.yaml"

# Start agent
ExecStart=/home/unitree/.local/bin/uv run /home/unitree/roboai-espeak/src/run.py lex_channel_chief

# Restart on failure (with delay to avoid rapid restart loops)
Restart=on-failure
RestartSec=15
StartLimitInterval=300
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lex_agent

[Install]
WantedBy=multi-user.target
