# robostore-ai.service
[Unit]
Description=Robostore AI Web Interface
After=network.target multi-user.target
Before=astra-vein-autostart.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/roboai
ExecStart=/usr/local/bin/uv run python robostore_ai/server.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=false
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
